# etc/supervisor/conf.d/php-fpm.conf
[program:php-fpm]
command = /usr/sbin/php-fpm7.3 -F -y /etc/php/7.3/fpm/php-fpm.conf
user = root
autostart = true
